# files/default/supervisor-ckan-harvest.conf
; ===================================================================
; Supervisor configuration for CKAN Harvest background fetch consumer
; ===================================================================

; 1. Copy this file to /etc/supervisord.d
; 2. Make sure the paths below match your setup


[program:ckan-harvest-fetch]

; Use the full paths to the virtualenv and your configuration file here.
command=/usr/lib/ckan/default/bin/ckan_cli harvester fetch_consumer

environment=PASTER_PLUGIN=ckanext-harvest

; User the fetch consumer runs as.
user=ckan


; Start just a single fetch consumer.
numprocs=1
process_name=%(program_name)s-%(process_num)02d


; Log files.
stdout_logfile=/var/log/ckan/ckan-harvest-fetch.log
stderr_logfile=/var/log/ckan/ckan-harvest-fetch.log


; Make sure that the fetch consumer is started on system start and automatically
; restarted if it crashes unexpectedly.
autostart=true
autorestart=true


; Number of seconds the process has to run before it is considered to have
; started successfully.
startsecs=10

; ====================================================================
; Supervisor configuration for CKAN Harvest background gather consumer
; ====================================================================

; 1. Copy this file to /etc/supervisord.d
; 2. Make sure the paths below match your setup


[program:ckan-harvest-gather]

; Use the full paths to the virtualenv and your configuration file here.
command=/usr/lib/ckan/default/bin/ckan_cli harvester gather_consumer

environment=PASTER_PLUGIN=ckanext-harvest


; User the gather consumer runs as.
user=ckan


; Start just a single gather consumer.
numprocs=1
process_name=%(program_name)s-%(process_num)02d


; Log files.
stdout_logfile=/var/log/ckan/ckan-harvest-gather.log
stderr_logfile=/var/log/ckan/ckan-harvest-gather.log


; Make sure that the gather consumer is started on system start and automatically
; restarted if it crashes unexpectedly.
autostart=true
autorestart=true


; Number of seconds the process has to run before it is considered to have
; started successfully.
startsecs=10
